# files/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:init]
command=/root/vmail/init.sh
autorestart=false

[program:rsyslog]
command=/etc/init.d/rsyslog start
autorestart=false
